# watertank.service
[Unit]
Description=Water Tank Sensor Flask App
After=network.target

[Service]
User=admin
WorkingDirectory=/home/admin/watertankvolume
ExecStart=/home/admin/watertankvolume/venv/bin/python /home/admin/watertankvolume/watertankvolume.py
Restart=always
Environment="PATH=/home/admin/watertankvolume/venv/bin"

[Install]
WantedBy=multi-user.target
